# other/crontab
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.rYN5p3qCJ7 installed on Sun Mar 10 23:27:33 2013)
# (Cron version V5.0 -- $Id: crontab.c,v 1.12 2004/01/23 18:56:42 vixie Exp $)
0,30 * * * *  /cygdrive/c/Program\ Files\ \(x86\)/Growl\ for\ Windows/growlnotify.exe /t:"now" "`date -R`" 

